# config/zfs/syncoid.timer
[Unit]
Description=zfast,zslow Backup Firing

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
